# etc/profile.d/05env.shrc
EDITOR=vim
export EDITOR

PATH=$HOME/bin:$HOME/.local/bin:$PATH
export PATH
