# etc/profile.d/10alias.shrc
alias grep='grep --color=auto'
alias egrep='egrep --color=auto'
alias ls='ls --color=auto'
alias mv='mv -i'
alias rm='rm -i'
alias cp='cp -i'
alias screen='screen -a'
alias ltree='tree -AFC'
alias ll='ls -lah'
