# etc/profile.d/20alias.zshrc
# disable correction for some commands
alias mv='nocorrect mv -i'
alias cp='nocorrect cp -i'
alias mkdir='nocorrect mkdir'
